# net-mgmt/collectd5/files/openrc-collectd.in
#!/sbin/openrc-run

# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# collectd_flags (string):      Set to "" by default.
#               Additional flags to collectd, see manual page.

depends()
{
	provide collectd
	after ntpd dbus
}

name="collectd"
command=%%PREFIX%%/bin/${name}
command_args="-f ${collectd_flags}"
pidfile=/var/run/collectd.pid
supervisor=supervise-daemon
output_log="/var/log/collectd.log"
error_log="/var/log/collectd.err"
