# trident/trident-core/pkg-plist
/boot/loader.conf.local
/etc/rc.conf.d/trident.conf
/etc/sysctl.d/trident-core.conf
bin/check-tor-mode
bin/disable-tor-mode
bin/enable-tor-mode
bin/mk_zfs_homedir
etc/devd/trident-media.conf
etc/fonts/local.conf.trident
etc/init.d/net-online
etc/init.d/trident_init
etc/openal/alsoft.conf
etc/sudoers.d/trident
etc/trident/ipfw.tor.rules
etc/trident/recommended-services
etc/trident/required-services
etc/trident/torrc
etc/trueos.conf.dist
etc/trueos/recommended-services-core
etc/trueos/recommended-services-desktop
etc/trueos/required-services-core
etc/trueos/required-services-desktop
share/cups/data/default-testpage.orig.pdf
share/fluxbox/styles/bora_blue_4k/theme.cfg
share/icons/default/index.theme
share/icons/hicolor/256x256/places/start-here-lumina.png
share/trident/scripts/detect-best-driver.sh
share/trident/scripts/generate-desktop-entry
share/trident/scripts/generate-xorg-conf.sh
share/trident/scripts/remove-desktop-entry
share/trident/scripts/reset-firewall
share/trident/scripts/reset-firewall-tor
share/trident/scripts/sys-init.sh
share/trident/scripts/validate-services.sh
share/trident/xorg-templates/xorg.conf
share/trident/xstartup-parallel/disable-beep.sh
share/trident/xstartup-parallel/load-xresources.sh
share/trident/xstartup/01_xdg-user-dirs.sh
share/trident/xstartup/02_fix-gtk-environ.sh
share/trident/xstartup/enable-inputmethod.sh
share/trident/xstartup/gpg-agent.sh
share/trident/xstartup/ssh-agent.sh
share/trueos/scripts/mkzfsdir.sh
share/trueos/scripts/rmzfsdir.sh
