# sysutils/sysup/files/openrc-sysup.in
#!/sbin/openrc-run
#

name=sysup
command="%%PREFIX%%/bin/sysup"
pidfile=/var/run/sysup.pid
command_args="-websocket"
supervisor=supervise-daemon

depend()
{
	need localmount
	use net
}
