# sysutils/sysup/files/sysup.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: sysup
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# consul_sysup (bool):	Set to NO by default.
#				Set it to YES to enable consul.

. /etc/rc.subr

name=sysup
rcvar=sysup_enable

load_rc_config $name

: ${sysup_enable:="NO"}

pidfile=/var/run/sysup.pid
procname="%%PREFIX%%/bin/sysup"
command="/usr/sbin/daemon"
command_args="-f -p ${pidfile} /usr/bin/env ${consul_env} ${procname} -websocket"

run_rc_command "$1"
