# x11/lightdm/files/openrc-lightdm.in
#!/sbin/openrc-run
# OpenRC script for lightdm
# Copyright 2017 by Ken Moore <ken@ixsystems.com>
# Available under 2-clause BSD license
# Written May 9th, 2017

name="lightdm"
description="The Simple Login Manager (SLiM)"
command=/usr/local/sbin/lightdm
command_args=""
start_stop_daemon_args="--background"
pidfile="/var/run/lightdm.pid"

depend(){
        provide xlogin
        keyword -login -shutdown
        need localmount dbus hald
}
