# sysutils/ec2-scripts/files/ec2_configinit.in
#!/sbin/openrc-run

# Service dependencies
depend()
{
        need localmount
        after network
	before ec2_fetchkey
} 

name="ec2_configinit"

CONFIGURL="http://169.254.169.254/latest/user-data"

start()
{
	# Download to a temporary location.
	einfo "Fetching EC2 user-data"
	CONFFILE=`mktemp "${TMPDIR:-/tmp}/configinit.XXXXXX"`
	fetch -o ${CONFFILE} ${CONFIGURL} 2>/dev/null

	# If we succeeded, process it; otherwise report failure.
	if [ $? -eq 0 ]; then
		# Process the user-data.
		echo .
		einfo "Processing EC2 user-data"
		/usr/local/sbin/configinit $CONFFILE
		echo .
	else
		ewarn " failed."
	fi

	# Whether we suceeded or not, delete the temporary file.
	rm $CONFFILE

	# Signal /etc/rc to reload rc.conf in case it changed.
	kill -ALRM $$
}

stop()
{
	return 0
}
