# deskutils/lazylibrarian/files/openrc-lazylibrarian.in
#!/sbin/openrc-run
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# %%PORTNAME%%_user:  The user account %%PORTNAME%% daemon runs as what
#     you want it to be. It uses 'media' user by
#     default. Do not sets it as empty or it will run
#     as root.
# %%PORTNAME%%_group: The group account %%PORTNAME%% daemon runs as what
#     you want it to be. It uses 'media' group by
#     default. Do not sets it as empty or it will run
#     as wheel.
# %%PORTNAME%%_data_dir:  Directory where %%PORTNAME%% configuration
#     data is stored.
#     Default: /var/db/%%PORTNAME%%

name=lazylibrarian

: ${%%PORTNAME%%_user:="librarian"}
: ${%%PORTNAME%%_group:="librarian"}
: ${%%PORTNAME%%_data_dir:="/var/db/%%PORTNAME%%"}

pidfile="/var/run/%%PORTNAME%%/%%PORTNAME%%.pid"
command="%%PYTHON_CMD%%"
command_args="%%DATADIR%%/LazyLibrarian/LazyLibrarian.py --datadir ${lazylibrarian_data_dir} --pidfile ${pidfile} --daemon --nolaunch"

start_pre() {
  if [ ! -d ${pidfile%/*} ]; then
    install -d -o ${%%PORTNAME%%_user} -g ${%%PORTNAME%%_group} ${pidfile%/*}
  fi

  if [ ! -d ${%%PORTNAME%%_data_dir} ]; then
    install -d -o ${%%PORTNAME%%_user} -g ${%%PORTNAME%%_group} ${%%PORTNAME%%_data_dir}
  fi
}
