# x11/sddm/files/openrc-sddm.in
#!/sbin/openrc-run

name="sddm"
command=/usr/local/bin/sddm
start_stop_daemon_args="--background"
pidfile="/var/run/sddm.pid"

depend(){
        provide xlogin
        keyword -login -shutdown
        after localmount dbus moused
}

stop_pre(){
	pkill -u sddm || true
	killall sddm || true
}

stop(){
	return 0
}

stop_post(){
	killall X || true
}
