# net/mDNSResponder/files/openrc-mdnsresponderposix.in
#!/sbin/openrc-run

name=mdnsresponderposix

: ${mdnsresponderposix_pidfile="/var/run/${name}.pid"}

command="%%PREFIX%%/bin/mDNSResponderPosix"
hostname=`sysrc -n hostname`
command_args="-b -P ${mdnsresponderposix_pidfile} -n ${hostname}"

depend() {
	keyword -shutdown
	after network
	need mdnsd
}
