# x11/gdm/files/openrc-gdm.in
#!/sbin/openrc-run

name="gdm"
command="/usr/local/sbin/${name}"
pidfile="/var/run/${name}.pid"
start_stop_daemon_args="-m --background"

depend() {
        provide gdb
        need moused hald dbus syscons
}

stop_pre()
{
    pkill -9 -u gdm
    pkill -9 X
    pkill -9 gdm-session-worker
}

stop_post()
{
    killall Xorg
}
