# x11/gdm/files/openrc-gdm.in
#!/sbin/openrc-run

name="gdm"
command="/usr/local/sbin/${name}"
pidfile="/var/run/${name}.pid"
start_stop_daemon_args="-m --background"

depend() {
        provide gdb
        need moused hald dbus syscons
}

stop_pre()
{
    pkill -9 -f gdm >/dev/null 2>/dev/null || true
}

stop()
{
    return 0
}
